# neo/supervisord.conf
[supervisord]
nodaemon=true

[program:dotnet]
command=dotnet neo-cli.dll /rpc
